# myapp.service
[Unit]
Description=My App

[Service]
ExecStart=/home/vagrant/main
Type=simple

[Install]
WantedBy=multi-user.target